# mpc/rc.local
#!/bin/sh -e

# rc.local for Spider MCCs; this file is automatically executed at boot time.

# create mce devices ...?
/sbin/udevadm trigger --verbose

# find a /data drive
if [ -d /data3/mas ]; then
  ln -sf /data3 /data
elif [ -d /data2/mas ]; then
  ln -sf /data2 /data
elif [ -d /data1/mas ]; then
  ln -sf /data1 /data
else
  ln -sf /data0 /data
fi

/data/mas/bin/mpc_inf_loop &

exit 0
